# packaging/launchagent/consulting.botanica.sod.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Reference copy of the per-user LaunchAgent. The canonical template is embedded in
     SodKit/LaunchAgent.swift and written (with paths substituted) by `sod install`.
     Placeholders __SOD_BIN__ / __SOCK__ below are filled in at install time. It MUST be
     a LaunchAgent (GUI session), never a LaunchDaemon, or Touch ID cannot present. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>consulting.botanica.sod.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__SOD_BIN__</string>
        <string>ssh-agent</string>
        <string>-d</string>
        <string>-a</string>
        <string>__SOCK__</string>
    </array>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ProcessType</key><string>Interactive</string>
    <key>StandardOutPath</key><string>__SOCK__.log</string>
    <key>StandardErrorPath</key><string>__SOCK__.log</string>
</dict>
</plist>
